# ansible/machines/servers/metal/tardis/roles/configure-remote-backups/files/systemd/kopia-remote-file-backup.timer
[Unit]
Description=Run Kopia Remote File Backup Daily

[Timer]
OnCalendar=*-*-* 8:00:00
Persistent=true

[Install]
WantedBy=timers.target
